# phd.conf
[supervisord]
nodaemon=true

[program:PhabricatorRepositoryPullLocalDaemon]
command=/phab/phabricator/scripts/daemon/phd-daemon PhabricatorRepositoryPullLocalDaemon --phd=/tmp/phd/pid
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PhabricatorRepositoryPullLocalDaemon.log
stderr_logfile=/var/log/supervisor/PhabricatorRepositoryPullLocalDaemon_err.log

[program:PhabricatorGarbageCollectorDaemon]
command=/phab/phabricator/scripts/daemon/phd-daemon PhabricatorGarbageCollectorDaemon --phd=/tmp/phd/pid
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon.log
stderr_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon_err.log

[program:PhabricatorTaskmasterDaemon1]
command=/phab/phabricator/scripts/daemon/phd-daemon PhabricatorTaskmasterDaemon --phd=/tmp/phd/pid
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon1.log
stderr_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon1_err.log

[program:PhabricatorTaskmasterDaemon2]
command=/phab/phabricator/scripts/daemon/phd-daemon PhabricatorTaskmasterDaemon --phd=/tmp/phd/pid
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon2.log
stderr_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon2_err.log

[program:PhabricatorTaskmasterDaemon3]
command=/phab/phabricator/scripts/daemon/phd-daemon PhabricatorTaskmasterDaemon --phd=/tmp/phd/pid
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon3.log
stderr_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon3_err.log

[program:PhabricatorTaskmasterDaemon4]
command=/phab/phabricator/scripts/daemon/phd-daemon PhabricatorTaskmasterDaemon --phd=/tmp/phd/pid
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon4.log
stderr_logfile=/var/log/supervisor/PhabricatorGarbageCollectorDaemon4_err.log

[program:PhpFPM]
command=/usr/sbin/php5-fpm
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/php5_fpm.log
stderr_logfile=/var/log/supervisor/php5_fpm_err.log
